# vendor/golang.org/x/build/env/linux-x86-std/scripts/rc.local
#!/bin/sh -e

start-stop-daemon --start --exec=/usr/local/bin/stage0 --background /usr/local/bin/stage0 -- --network-wait=15s
exit 0
