# vendor/golang.org/x/build/env/linux-arm64/linaro/buildlet.service
[Unit]
Description=Go builder buildlet
After=network.target

[Install]
WantedBy=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'cd /usr/local/bin; /usr/bin/curl -R -f -z buildlet-stage0 -o buildlet-stage0 https://storage.googleapis.com/go-builder-data/buildlet-stage0.linux-arm64 && chmod +x buildlet-stage0'
ExecStart=/usr/local/bin/buildlet-stage0
Restart=always
RestartSec=2
StartLimitInterval=0
